# apache-airflow/systemd/airflow-worker.service
[Unit]
Description=Airflow celery worker daemon
After=network.target postgresql.service mysql.service redis.service rabbitmq-server.service airflow-flower.service
Wants=postgresql.service mysql.service redis.service rabbitmq-server.service

[Service]
User=ec2-user
Group=ec2-user
Type=simple
PermissionsStartOnly=True
RuntimeDirectory=airflow
RuntimeDirectoryMode=0775
PIDFile=/run/airflow/worker.pid
Environment=AIRFLOW_HOME=/home/ec2-user/airflow
ExecStart=/bin/bash -c 'source /home/ec2-user/anaconda3/bin/activate airflow && /home/ec2-user/anaconda3/envs/airflow/bin/airflow worker --pid /run/airflow/worker.pid'
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
